# usr/lib/systemd/system/greenboot-success.target
[Unit]
Description=GreenBoot Healthcheck Success Target
After=greenboot-healthcheck.service
Requires=boot-complete.target
