# usr/lib/systemd/system/greenboot-success.target
[Unit]
Description=GreenBoot Healthcheck Success Target
Requires=boot-complete.target
